# systemd/ssh-telegram-monitor.service
[Unit]
Description=SSH Telegram Monitor for Failed Login Attempts
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/ssh_monitor_advanced.py
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target